# repository_downloader/crontab
*/20 * * * * timeout 1200 python3 /app/script.py >> /var/log/cron.log 2>&1
